# images/playbook/roles/docker/files/ecr.timer
[Unit]
Description=Run ECR hourly

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
